# libvirt/binary/libirt-cloudpc-4.4.0/release/lib/systemd/system/virtlogd.socket
[Unit]
Description=Virtual machine log manager socket
Before=libvirtd.service

[Socket]
ListenStream=/usr/local/var/run/libvirt/virtlogd-sock

[Install]
WantedBy=sockets.target
